# deploy/systemd/dotmac-static-sync.service
[Unit]
Description=DotMac ERP - sync static assets to the Nginx web root
Documentation=file:/root/dotmac/scripts/sync-static.sh
After=network.target
# No-op unless the project (and its sync script) is present.
ConditionPathExists=/root/dotmac/scripts/sync-static.sh

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/root/dotmac
ExecStart=/root/dotmac/scripts/sync-static.sh
# Surface runs in the journal and a dedicated log alongside the app logs.
StandardOutput=append:/var/log/dotmac/static-sync.log
StandardError=append:/var/log/dotmac/static-sync.log
