# deploy/systemd/dotmac-static-sync.timer
[Unit]
Description=Periodically sync DotMac ERP static assets to the Nginx web root
Documentation=file:/root/dotmac/deploy/systemd/dotmac-static-sync.service

[Timer]
# Catch drift from any update path (bare `docker restart`, live bind-mount
# edits) shortly after boot, then every few minutes thereafter. A couple of
# minutes of lag is harmless; a missed sync is what caused stale JS to be
# served under `immutable` cache headers for up to 30 days.
OnBootSec=1min
OnUnitActiveSec=2min
AccuracySec=30s
# Run a missed sync after downtime instead of waiting for the next interval.
Persistent=true

[Install]
WantedBy=timers.target
